# daemon/fm-daemon.service
[Unit]
Description=FM Radio Daemon
After=network.target

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/fm-daemon.py
Restart=on-failure
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal
# Volume boost applied by ffmpeg. Use ffmpeg syntax: e.g. 6dB, 10dB, 2.0 (linear).
Environment=VOLUME_BOOST=6dB

[Install]
WantedBy=multi-user.target
